# etc/supervisord.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx 
redirect_stderr = true
autostart=true
autorestart=true
startretries=5

