# etc/supervisord.d/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm -F -O
redirect_stderr = true
autostart=true
autorestart=true
startretries=5
